# ubuntu-files/config/rc.local
#!/bin/bash
cd /root/icubtech/python-cameras/ubuntu-files/yarp.local
./preliminary.sh
exit 0
